# ansible/roles/monitoring/templates/self-hosted-healthcheck.service.j2
[Unit]
Description=Capture lightweight host health data for the private homelab

[Service]
Type=oneshot
ExecStart=/usr/local/bin/self-hosted-healthcheck
