# ansible/roles/backup/templates/self-hosted-backup.service.j2
[Unit]
Description=Local backup for self-hosted homelab data

[Service]
Type=oneshot
ExecStart=/usr/local/bin/self-hosted-backup
