# ansible/roles/backup/templates/self-hosted-backup.timer.j2
[Unit]
Description=Run local self-hosted backup schedule

[Timer]
OnCalendar={{ platform_backup_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
